# deploy/systemd/ollama_system.service

[Unit]
Description=ollama_system FastAPI service
After=network.target

[Service]
User=bradi
Group=bradi
WorkingDirectory=/home/bradi/ollama_system
Environment="PATH=/home/bradi/ollama_system/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Run a single-process Uvicorn to reduce memory usage on low-RAM devices
ExecStart=/home/bradi/ollama_system/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
